# hortlak.conf
[program:hortlak]
command="/usr/bin/mono-service --no-daemon -d:/app /app/Hortlak.exe"
